# launchagents/ai.openclaw.taskdaemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.taskdaemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>--directory</string>
        <string>/home/hurin/.openclaw/monitor</string>
        <string>python</string>
        <string>/home/hurin/.openclaw/monitor/task-daemon.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/home/hurin/.openclaw/monitor</string>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/home/hurin/.openclaw/monitor/daemon-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/home/hurin/.openclaw/monitor/daemon-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/home/hurin/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/hurin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
